# backend/systemd/m5stack-ha-publish-chat.timer
[Unit]
Description=Run M5Stack LLM chat smoke MQTT publisher every 10 minutes

[Timer]
OnBootSec=180s
OnUnitActiveSec=10min
AccuracySec=30s
Unit=m5stack-ha-publish-chat.service
Persistent=true

[Install]
WantedBy=timers.target
