# backend/systemd/m5stack-ha-publish.service
[Unit]
Description=Publish M5Stack LLM health to Home Assistant MQTT
After=network-online.target llm-openai-api.service llm-sys.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/m5stack-ha-publish
