# backend/systemd/power-sentinel-stackflow-unit.service
[Unit]
Description=Power Sentinel StackFlow custom unit
Documentation=https://github.com/m5stack/StackFlow
After=llm-sys.service power-sentinel-api.service
Wants=llm-sys.service power-sentinel-api.service

[Service]
Type=simple
ExecStart=/usr/local/bin/power-sentinel-stackflow-unit --unit sentinel --api-url 'http://127.0.0.1:8088/api/v1/summary?stackflow_safe=1' --timeout 4
Restart=on-failure
RestartSec=3s
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=full
# Must remain false/default: this service has to bind /tmp/rpc.sentinel and
# reach llm_sys callback sockets under /tmp/llm/*.sock.
PrivateTmp=false

[Install]
WantedBy=multi-user.target
