# backend/systemd/power-sentinel-serial-bridge.service
[Unit]
Description=Power Sentinel CoreS3 serial bridge
After=power-sentinel-api.service
Wants=power-sentinel-api.service

[Service]
Type=simple
Environment=POWER_SENTINEL_SERIAL_PORT=/dev/ttyS1
Environment=POWER_SENTINEL_SERIAL_BAUD=115200
Environment=POWER_SENTINEL_SUMMARY_URL=http://127.0.0.1:8088/api/v1/summary
ExecStart=/usr/local/bin/power-sentinel-serial-bridge
Restart=on-failure
RestartSec=2s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
